# scripts/com.remedymatch.ghl-write-drain.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.remedymatch.ghl-write-drain</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/remedymatch/deploy-chat/scripts/run_ghl_write_drain.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/remedymatch/Library/Logs/ghl-write-drain.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/remedymatch/Library/Logs/ghl-write-drain.err.log</string>
</dict>
</plist>
